# server/showtrack-sync.service
[Unit]
Description=ShowTrack sync server
After=network.target

[Service]
Type=simple
# EDIT these two lines to match your machine:
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/showtrack/server
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
# optional: change where data is stored
# Environment=DATA_DIR=/home/YOUR_USERNAME/showtrack-data
# Environment=PORT=8570

[Install]
WantedBy=multi-user.target
